# version1/config/ripitt.gunicorn.socket
Syntax:- 
[Unit]
Description=ripitt.gunicorn socket

[Socket]
ListenStream=/run/ripitt.gunicorn.sock

[Install]
WantedBy=sockets.target